# dist/trusted-network-gateway.service
[Unit]
Description=Trusted Network Gateway Daemon
Documentation=https://github.com/inclavare-containers/tng/blob/master/docs/configuration.md
After=network.target

[Service]
ExecStart=/usr/bin/tng launch --config-file /etc/tng/config.json
Restart=always

[Install]
WantedBy=multi-user.target
